# data_files/upstart/crewchief.conf
description "Rackconnect Crewchief"
start on mounted MOUNTPOINT=/proc/xen TYPE=xenfs
task
exec crewchief
