# data_files/upstart/crewchief.conf
description "Rackconnect Crewchief"
start on started rc
task
env COUNT=10
env INTERVAL=60
exec crewchief --count ${COUNT} --interval ${INTERVAL}
