# data_files/upstart/crewchief.conf
description "Rackconnect Crewchief"
start on runlevel [2345]
task
exec crewchief
