# data_files/upstart/crewchief.conf
description "Rackconnect Crewchief"
start on started rc
task
exec crewchief
