# data_files/upstart-ubuntu/rcrun.conf
description "rcrun"
start on started rc-sysinit
task
env COUNT=10
env INTERVAL=60
exec rcrun --count ${COUNT} --interval ${INTERVAL}
